# dot_config/systemd/user/serena-mcp.service
[Unit]
Description=Serena MCP server for Codex

[Service]
Type=simple
ExecStart=/home/birbante/.local/bin/serena start-mcp-server --context=codex --transport streamable-http --host 127.0.0.1 --port 9121 --enable-web-dashboard true --open-web-dashboard true
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
